# ee/profile.d/00-xdg-docker-config.sh
# shellcheck shell=sh
export XDG_CACHE_HOME="${XDG_CACHE_HOME:-/var/run/user/$(id -u)/cache}"