# ee/profile.d/p10k-centos.zsh
_p9k__force_must_init=1
export CENTOS_POWERLEVEL9K_CONFIG_FILE=${CENTOS_POWERLEVEL9K_CONFIG_FILE:-/etc/profile.d/p10k-centos.zsh}
# Closest match to: #bc1142  = #d7005f
# Closest match to: #a14f8c  = #af5f87
# Closest match to: #f0a91f  = #ffaf00
typeset -g POWERLEVEL9K_OS_ICON_FOREGROUND=214 # #ffaf00 orange1
typeset -g POWERLEVEL9K_OS_ICON_BACKGROUND=161 # #d7005f magenta5
typeset -g POWERLEVEL9K_OS_ICON_BACKGROUND=132 # #af5f87 magenta4
# Custom icon.
typeset -g POWERLEVEL9K_OS_ICON_CONTENT_EXPANSION=''

export POWERLEVEL9K_OS_ICON_FOREGROUND POWERLEVEL9K_OS_ICON_BACKGROUND POWERLEVEL9K_OS_ICON_CONTENT_EXPANSION
