# ee/profile.d/p10k-centos.zsh
# Closest match to: #bc1142  = #a14f8c
# Closest match to: #f0a91f  = #e6cd69
typeset -g POWERLEVEL9K_OS_ICON_FOREGROUND=003 # #e6cd69 yellow
typeset -g POWERLEVEL9K_OS_ICON_BACKGROUND=005 # #a074c4 magenta
# Custom icon.
typeset -g POWERLEVEL9K_OS_ICON_CONTENT_EXPANSION=''
